# systemd/system/ewwpoweroff.service
[Unit]
Description=Prevent arek arek to touch this thing.

[Service]
ExecStart=/usr/bin/poweroff

[Install]
Alias=ewwpoweroff.service
